# networkMonitor/networkMonitor.service
[Unit]
Description=Notify of external ip address change

[Service]
Type=simple
ExecStart=/etc/systemd/system/networkMonitor.script
